# app/ops/morchess.service
[Unit]
Description=Mor Chess 3 API
After=network.target

[Service]
Type=simple
User=morchess
WorkingDirectory=/var/www/morchess-api
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5

# Safety
Environment=NODE_ENV=production
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target